# deploy/grayboard.service
[Unit]
Description=grayboard server
After=network.target

[Service]
Type=simple
User=grayboard
WorkingDirectory=/opt/grayboard
EnvironmentFile=/etc/grayboard/env
ExecStart=/usr/local/bin/bun run src/server/main.ts
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
